# my-default-settings/files/etc/hotplug.d/ieee80211/99-mywifi
#!/bin/sh

[ "${ACTION}" = "add" ] || exit 0
[ -z "$(uci -q get base_config.@status[0].wifi)" ] && . /sbin/mywifi &
